# docker/rootfs/etc/auto-scan/smanga_daemon.ini
[program:smanga-daemon]
command=php7 artisan daemon:start
directory=/app/php
user=smanga
stopwaitsecs=600
autorestart=true
stdout_logfile=/logs/smanga-daemon.log
redirect_stderr=true
