# docker/rootfs/etc/auto-scan/smanga_compress.ini
[program:smanga-compress]
command=php7 artisan queue:listen --queue=compress
directory=/app/php
numprocs=2
process_name=%(program_name)s_%(process_num)02d
user=smanga
autorestart=true
stdout_logfile=/logs/smanga-compress.log
redirect_stderr=true
