# docker/rootfs/etc/auto-scan/smanga_scan.ini
[program:smanga-scan]
command=php7 artisan queue:listen --queue=scan
directory=/app/php
numprocs=1
user=smanga
autorestart=true
stdout_logfile=/logs/smanga-scan.log
redirect_stderr=true
