# docker/rootfs/etc/supervisor.d/smanga_websocket.ini
[program:smanga-websocket]
command=php7 artisan socket start
directory=/app/php
user=smanga
stopwaitsecs=3600
autorestart=true
stdout_logfile=/data/logs/smanga-websocket.log
redirect_stderr=true
