# docker/rootfs/etc/supervisor.d/smanga_scan.ini
[program:smanga-queue]
command=php7 artisan queue:work --queue=delete,scan
directory=/app/php
numprocs=1
stopwaitsecs=3600
user=smanga
autorestart=true
stdout_logfile=/data/logs/smanga-scan.log
redirect_stderr=true
