# docker/rootfs/etc/supervisor.d/smanga_compress.ini
[program:smanga-queue]
command=php7 artisan queue:work --queue=compress
directory=/app/php
numprocs=1
stopwaitsecs=3600
user=smanga
autorestart=true
stdout_logfile=/data/logs/smanga-compress.log
redirect_stderr=true
