# docker/rootfs/etc/supervisor.d/smanga_daemon.ini
[program:smanga-daemon]
command=php7 artisan daemon:start
directory=/app/php
user=smanga
stopwaitsecs=600
autorestart=true
stdout_logfile=/data/logs/smanga-daemon.log
redirect_stderr=true
