# docker/alpine/rootfs/etc/services.d/nginx/finish
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

/usr/sbin/nginx -s quit