# docker/rootfs/etc/services.d/smanga_compress/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec s6-setuidgid smanga \
    php7 /app/php/laravel/artisan queue:listen --queue=compress