# docker/rootfs/etc/s6-overlay/s6-rc.d/init-030-config/up
/etc/s6-overlay/s6-rc.d/init-030-config/run