# docker/rootfs/etc/s6-overlay/s6-rc.d/init-021-touchfiles/up
/etc/s6-overlay/s6-rc.d/init-021-touchfiles/run