# docker.debian/rootfs/etc/services.d/log/run
#!/usr/bin/with-contenv bash

sleep 4

exec s6-setuidgid root tail -f /log/*.log