# docker/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid root /usr/sbin/nginx -g "daemon off;"