# docker/rootfs/etc/s6-overlay/s6-rc.d/svc-logs_out/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec s6-setuidgid root \
    tail -f ${OUT_LOG_FILES}