# docker/rootfs_mysql/etc/s6-overlay/s6-rc.d/init-120-mysql-fixuser/up
/etc/s6-overlay/s6-rc.d/init-120-mysql-fixuser/run