# docker/alpine/rootfs/etc/services.d/php7/finish
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

kill -s QUIT $(ps aux | grep php-fpm | grep -v grep | awk '{print $1}')