# docker.debian/rootfs/etc/services.d/php7/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid root /usr/sbin/php-fpm7.4 -F -R