# docker/rootfs/auto-scan/smanga_compress.ini
[program:smanga_compress]
command=php7 artisan queue:listen --queue=compress
directory=/app/php
user=smanga
autorestart=true
stdout_logfile=/logs/smanga_compress.log
redirect_stderr=true
