# docker/rootfs/auto-scan/smanga_scan.ini
[program:smanga_scan]
command=php7 artisan queue:listen --queue=scan
directory=/app/php
user=smanga
autorestart=true
stdout_logfile=/logs/smanga_scan.log
redirect_stderr=true
