# docker/rootfs/auto-scan/smanga_websocket.ini
[program:smanga_websocket]
command=php7 artisan socket start
directory=/app/php
user=smanga
autorestart=true
stdout_logfile=/logs/smanga_websocket.log
redirect_stderr=true
