# docker/rootfs/etc/services.d/logs_out/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec s6-setuidgid root \
    tail -q -f /logs/*.log | \
    grep -vFf <(find "${EXCLUDE_LOG}" -type f -name "*.log" -print0 | xargs -0 cat)