# docker/rootfs/etc/services.d/logs_out/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec s6-setuidgid root tail -f /logs/nginx_error.log